# .chef/cookbooks/pincrowd/files/default/etc/profile.d/php-composer.sh
#!/bin/sh
alias composer='php /usr/local/bin/composer.phar'
